# bgb.conf
[program:bgb]
command=wine /root/bgb.exe
autorestart=false